# root/etc/services.d/urserver/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /app/urserver/urserver \
	--debug --daemon --pidfile /config/urserver.pid
